# systemd/fan-control.service
[Unit]
Description=Fan control
After=syslog.target

[Service]
Type=simple
WorkingDirectory=/opt/fan-control
ExecStart=/opt/fan-control/fan-control.py --device C2 --update-delay 5 --quiet
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
